# deliverables/build-weekly-synthesis-report/com.chipai.mission-control.weekly-synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.chipai.mission-control.weekly-synthesis</string>
	
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>cd ~/Documents/Shared/projects/build-weekly-synthesis-report && /usr/bin/node scripts/weekly-synthesis.js</string>
	</array>
	
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>5</integer>
		<key>Hour</key>
		<integer>16</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	
	<key>StandardOutPath</key>
	<string>/Users/chipai/.mission-control/weekly-synthesis.log</string>
	
	<key>StandardErrorPath</key>
	<string>/Users/chipai/.mission-control/weekly-synthesis-error.log</string>
	
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/chipai</string>
	</dict>
	
	<key>RunAtLoad</key>
	<false/>
	
	<key>KeepAlive</key>
	<false/>
</dict>
</plist>
